# p11-kit-trigger.service
[Unit]
Description=Rebuild certificate database
BindsTo=update-triggers.target

[Service]
Type=oneshot
ExecStart=/usr/bin/trust-certs


